# farm-island/scripts/farm-deploy.service
[Unit]
Description=SCM Lab Farm Island GitOps Deploy
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/scm-lab/farm-island/scripts/deploy.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=farm-deploy
