# farm-island/scripts/farm-backup.service
# =============================================================
# systemd Service Unit — Farm Island Backup
# =============================================================
# Installation:
#   sudo cp scripts/farm-backup.service /etc/systemd/system/
#   sudo cp scripts/farm-backup.timer   /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now farm-backup.timer
#
# Run manually:
#   sudo systemctl start farm-backup.service
#
# Check status:
#   systemctl status farm-backup.timer
#   journalctl -u farm-backup.service -n 50
# =============================================================

[Unit]
Description=Farm Island Backup (Docker Volumes + Databases)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=farm
WorkingDirectory=/home/farm/farm-island

# Backup destination (local)
Environment=BACKUP_ROOT=/home/farm/farm-backups

# Optional sync to network share — uncomment and adjust:
# Environment=REMOTE_DEST=farm@nas.local:/backups/farm-island

ExecStart=/home/farm/farm-island/scripts/backup.sh

# Log to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=farm-backup

# Timeout: 30 minutes (increase for large ERPNext databases if needed)
TimeoutStartSec=1800
