# farm-insel/scripts/farm-deploy.service
[Unit]
Description=SCM Labor Farm-Insel GitOps Deploy
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/opt/scm-labor/farm-insel/scripts/deploy.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=farm-deploy
